# deploy/st-albans.service
[Unit]
Description=St Albans Bin Collection App
After=network-online.target
Wants=network-online.target
# Wait a bit after network is available for DNS to be ready
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=zander
Group=zander
WorkingDirectory=/opt/st-albans/server
# Optional .env file; the leading '-' means ignore if missing
EnvironmentFile=-/opt/st-albans/server/.env
Environment=NODE_ENV=production
# Explicit PATH so the nvm Node is found without shell init
Environment=PATH=/home/zander/.nvm/versions/node/v22.20.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Use absolute Node path (update if Node version changes)
ExecStart=/home/zander/.nvm/versions/node/v22.20.0/bin/node dist/server.js
Restart=always
RestartSec=10
# Restart even on successful exit (in case of network issues)
RestartForceExitStatus=0

# Baseline hardening (comment out if troubleshooting):
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
ReadWritePaths=/opt/st-albans/server

[Install]
WantedBy=multi-user.target