# deploy/st-albans.service
[Unit]
Description=St Albans Bin Collection App
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/st-albans/server
EnvironmentFile=/opt/st-albans/server/.env
ExecStart=/usr/bin/env node dist/server.js
Restart=on-failure
RestartSec=5
User=pi

[Install]
WantedBy=multi-user.target
